# java8-ansible/ansible/roles/gradle/templates/profile.d/gradle.sh
export PATH=$PATH:{{gradle_base_dir}}/gradle-{{gradle_version}}/bin
